# ops/little-cloud.service
[Unit]
Description=Little Cloud and Magic Star FastAPI service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www
Group=www
WorkingDirectory=/www/wwwroot/little-cloud-and-magic-star/backend
EnvironmentFile=/www/wwwroot/little-cloud-and-magic-star/backend/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/www/wwwroot/little-cloud-and-magic-star/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8100 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=on-failure
RestartSec=3
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
